# files/powersave@.service
[Unit]
Description=Run the powersave script

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/powersave %I
