# files/powersave@.service
[Unit]
Description=Run the powersave script

[Service]
Type=oneshot
ExecStart=/usr/bin/powersave $I
