# rules.d/50-powersave-hdparm.rules
ACTION=="add", KERNEL=="sda", ATTR{queue/rotational}=="1", RUN+="/sbin/hdparm -S 4 -B 254 -a 2048 /dev/$kernel" 
