# rules/51-eth-down.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="eth*" RUN+="/sbin/ip link set dev %k down"
